# etc/systemd/user/dashcam.service
[Unit]
Description=Dash Camera
Wants=dashcam-cleaner.service

[Service]
ExecStart=/home/pi/dashcam/dashcam.py

[Install]
WantedBy=multi-user.target
