# etc/systemd/system/dashcam-cleaner.service
[Unit]
Description=Disk Cleaner for Dash Camera

[Service]
ExecStart=/home/pi/dashcam/clean.py

[Install]
WantedBy=multi-user.target
